# dnssec-sldns/etc/cron.d/bind-sign-zones
*/15 * * * *	root	/etc/bind/scripts/auto-signing.sh /etc/bind/zones 'dnsprovi.de dnssec.de'
