# dnssec-sldns/etc/cron.d/bind-sign-zones
@reboot		root	while ! dig . @localhost ; do sleep 10 ; done ; /etc/bind/scripts/auto-signing.sh /etc/bind/zones >> /var/log/named/auto-signing.log 2>&1
*/15 * * * *	root	/etc/bind/scripts/auto-signing.sh /etc/bind/zones >> /var/log/named/auto-signing.log 2>&1
