# dnssec-rootns-a/etc/supervisor/conf.d/dnssec-attendee.conf
[supervisord]
nodaemon=true

[program:postfix]
command=/etc/init.d/postfix start

[program:
command=/etc/init.d/cron start
