# dnssec-rootns/etc/cron.d/bind-sign-zones
@reboot 	root	( /etc/bind/scripts/sign-zone.sh . ; /etc/bind/scripts/sign-zone.sh test. ) >> /var/log/named/sign-zone.log 2>&1 ; while ! dig . @localhost ; do sleep 10 ; done ; rndc reload
*/15 * * * * 	root	( /etc/bind/scripts/sign-zone.sh . ; /etc/bind/scripts/sign-zone.sh test. ) >> /var/log/named/sign-zone.log 2>&1 ; rndc reload
