# dnssec-rootns-b/etc/supervisor/conf.d/dnssec-rootns-b.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:named]
command=/usr/sbin/named -u bind -f -4 -c /etc/bind/named.conf
