# dnssec-tldns/etc/cron.d/bind-zone-update
# Update records in TLD 
*/5 * * * *	root	/etc/bind/scripts/zone-update.sh >> /var/log/bind-zone-update.log
